# templates/apiserver-connectivity.timer.j2
[Timer]
OnActiveSec=0
OnUnitInactiveSec=10s

[Install]
WantedBy=timers.target
