# templates/stop-bird.service.j2
[Unit]
Description=Stops bird

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop bird.service
